# .config/systemd/user/openclaw-hk-tech-delivery.timer
[Unit]
Description=Run HK-Tech markdown delivery poller during 17:00-17:10

[Timer]
OnCalendar=*-*-* 17:00:00
OnCalendar=*-*-* 17:01:00
OnCalendar=*-*-* 17:02:00
OnCalendar=*-*-* 17:03:00
OnCalendar=*-*-* 17:04:00
OnCalendar=*-*-* 17:05:00
OnCalendar=*-*-* 17:06:00
OnCalendar=*-*-* 17:07:00
OnCalendar=*-*-* 17:08:00
OnCalendar=*-*-* 17:09:00
OnCalendar=*-*-* 17:10:00
Persistent=true
Unit=openclaw-hk-tech-delivery.service

[Install]
WantedBy=timers.target
